# preload-ng/preload-ng.service
[Unit]
Description=Preload-NG Daemon - Adaptive Readahead
Documentation=https://github.com/miguel-b-p/preload-ng
After=local-fs.target network.target

[Service]
Type=simple
# The binary is installed to /usr/sbin/preload. 
# We use --foreground for systemd and explicitly set the logfile to empty 
# so it logs to stderr, which systemd captures.
ExecStart=/usr/sbin/preload --foreground --logfile ""
# SIGHUP reloads the configuration
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
# Security hardening (daemon)
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
# StateDirectory automatically creates /var/lib/preload and handles permissions.
StateDirectory=preload
# Allow writing to the state directory and the default log directory
ReadWritePaths=/var/log
